# etc/cron.d/pwfcck
# MIN HOUR DAY MONTH DAYOFWEEK
# * * * * *
# Cronjob for security issues
00 04 * * * /usr/sbin/pwfcck
~
